# mother_supervisor.conf
[program:mother]
command=python /home/jdsnape/workers/mother/mother.py -c /home/jdsnape/workers/mother/mother.conf 
